# dockerconf/supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node /app/boot.js
autostart=true
autorestart=true
